# backend/deploy/medc-backend.service
## systemd unit for Gunicorn running the Flask app
[Unit]
Description=MedC Annotation Backend (Gunicorn)
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/droot/medc-img-annotation-app/backend
Environment="APP_ENV=production"
EnvironmentFile=/home/droot/medc-img-annotation-app/backend/.env
ExecStart=/home/droot/medc-img-annotation-app/backend/venv/bin/gunicorn -w 4 -k sync -b 0.0.0.0:5000 'app:create_app()'
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
User=droot
Group=droot
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target